# config/includes.chroot/etc/network/if-pre-up.d/iptables
#!/bin/sh
/sbin/iptables-restore < /etc/iptables.save
/sbin/ip6tables-restore < /etc/ip6tables.save
